# deploy/systemd/starmapd-deploy.service
[Unit]
Description=StarMap release deploy poller
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env bash __WORKDIR__/deploy/deploy-starmapd.sh --channel-url __CHANNEL_URL__ --service-name __SERVICE_NAME__ --health-url __HEALTH_URL__

